# updates for site/tc-enterprises-site-autoblog/nginx/tc-enterprises.service
[Unit]
Description=T&C Enterprises Next.js Application
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/var/www/tc-enterprises.com/updates-for-site/tc-enterprises-site-autoblog
Environment=NODE_ENV=production
Environment=PORT=3000
# Add your environment variables below (uncomment and fill in):
# Environment=SMTP_HOST=your-smtp-host
# Environment=SMTP_PORT=587
# Environment=SMTP_USER=your-email@domain.com
# Environment=SMTP_PASS=your-password
# Environment=SMTP_SECURE=false
# Environment=CONTACT_TO=contact@tc-enterprises.com
# Environment=CONTACT_FROM="T&C Website <contact@tc-enterprises.com>"
ExecStart=/usr/bin/npm start
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=tc-enterprises-app

# Security settings
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ReadWritePaths=/var/www/tc-enterprises.com
ProtectHome=yes

[Install]
WantedBy=multi-user.target