# tc-enterprises-site-autoblog/tc-enterprises-site-autoblog/nginx/tc-enterprises.service
[Unit]
Description=T&C Enterprises Next.js Application
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/var/www/tc-enterprises.com/current
Environment=NODE_ENV=production
Environment=PORT=3011
# SMTP Configuration - REQUIRED for contact form
# Uncomment and fill in with your actual values:
# Environment=SMTP_HOST=your-mail-server.com
# Environment=SMTP_PORT=587
# Environment=SMTP_USER=contact@tc-enterprises.com
# Environment=SMTP_PASS=your-email-password
# Environment=SMTP_SECURE=false
# Environment=SMTP_INSECURE=false
# Environment=CONTACT_TO=contact@tc-enterprises.com
# Environment=CONTACT_FROM="T&C Website <contact@tc-enterprises.com>"
ExecStart=/usr/bin/npm start
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Security settings
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ReadWritePaths=/var/www/tc-enterprises.com
ProtectHome=yes

[Install]
WantedBy=multi-user.target